# automation/ansible/ubuntu/TCGUI_install/tcgui.service
[Unit]
Description=tcgui start script

[Service]
ExecStart=/home/fortinet/tcgui.sh

[Install]
WantedBy=multi-user.target

